# infrastructure/systemd/tendermind-api.service
[Unit]
Description=TenderMind API
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/tender
ExecStart=/usr/bin/bash /root/tender/infrastructure/scripts/start-api.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
