# infrastructure/systemd/tendermind-web.service
[Unit]
Description=TenderMind Web
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/tender
Environment=PORT=3010
ExecStart=/usr/bin/bash /root/tender/infrastructure/scripts/start-web.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
